# user_systemd/.config/systemd/user/backup.timer
[Unit]
Description=Run Backup hourly

[Timer]
OnBootSec=1min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
