# user_systemd/.config/systemd/user/backup-prune.timer
[Unit]
Description=Run Backup Prune Daily

[Timer]
OnBootSec=10min
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
