# user_systemd/.config/systemd/user/dunst.service
[Unit]
Description=Dunst
After=xorg.target
Requires=xorg.target

[Service]
ExecStart=/usr/bin/dunst

[Install]
WantedBy=xorg.target
