# user_systemd/.config/systemd/user/xorg.target
[Unit]
Description=Xorg environment
BindsTo=graphical-session.target
